# ygege/rootfs/etc/s6-overlay/s6-rc.d/ygege/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Ygégé
# Runs the Ygégé application
# ==============================================================================

declare ygg_username
declare ygg_password
declare log_level

bashio::log.info "Starting Ygégé..."

# Get configuration from Home Assistant
ygg_username=$(bashio::config 'ygg_username')
ygg_password=$(bashio::config 'ygg_password')
log_level=$(bashio::config 'log_level')

# Validate configuration
if bashio::var.is_empty "${ygg_username}"; then
    bashio::log.fatal "YGG_USERNAME is not configured!"
    bashio::exit.nok "Please configure your YGG username in the addon configuration"
fi

if bashio::var.is_empty "${ygg_password}"; then
    bashio::log.fatal "YGG_PASSWORD is not configured!"
    bashio::exit.nok "Please configure your YGG password in the addon configuration"
fi

bashio::log.info "Configuration validated successfully"

# Export environment variables
export YGG_USERNAME="${ygg_username}"
export YGG_PASSWORD="${ygg_password}"
export LOG_LEVEL="${log_level}"
export BIND_IP="0.0.0.0"
export BIND_PORT="8715"

bashio::log.info "Starting Ygégé with user: ${ygg_username}"

# Run the application
exec /app/ygege
